# systemd/lte-collector.service
[Unit]
Description=LTE Remote Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/starlink_lte
ExecStart=/usr/bin/python3 /home/pi/starlink_lte/lte_remote_collector_en.py --serial-port auto --data-dir /home/pi/lte-collector-data --control-port 8897 --interval 0.5
Restart=always
RestartSec=1
StartLimitIntervalSec=0
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
