# systemd/lte-ground-station.service
[Unit]
Description=LTE Ground Station Receiver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/starlink_lte
ExecStart=/usr/bin/python3 /home/pi/starlink_lte/lte_ground_station.py --port 8079 --data-dir /opt/lte-ground-station-data
Restart=on-failure
RestartSec=2
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
